# sshd.sv.conf
[program:sshd]              
command=/etc/init.d/sshd start
autostart=true              
autorestart=true            
startsecs=3                 
